# arkdep-build.d/p14s/overlay/post_install/usr/lib/systemd/user/xdg-desktop-portal.service.d/sway.conf
[Unit]
PartOf=sway-session.target
After=sway-session.target
Requisite=sway-session.target
